# deploy/oracle/shorts-backend.service
[Unit]
Description=Shorts Publishing Studio Backend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/shorts-publishing-studio/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/shorts-publishing-studio/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
